# docs/services/vaultwarden/scripts/files/vaultwarden-update.timer
[Unit]
Description=Run Vaultwarden update Fridays at 04:00

[Timer]
OnCalendar=Fri *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target