# docs/services/vaultwarden/scripts/files/vaultwarden-backup.service
[Unit]
Description=Vaultwarden backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vaultwarden-backup.sh