# docs/services/vaultwarden/scripts/files/vaultwarden-update.service
[Unit]
Description=Vaultwarden weekly source-build update
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vaultwarden-update.sh